# templates/lita.service.erb
[Unit]
Description=Lita (%i) bot
After=network.target

[Service]
EnvironmentFile=-<%= @service_env_file %>-%i
Type=simple
ExecStart=<%= @bundler_binpath %>/bundle exec lita
ExecReload=<%= @bundler_binpath %>/bundle exec lita
ExecStop=<%= @bundler_binpath %>/bundle exec lita

[Install]
WantedBy=multi-user.target
